# supervisord.conf
[supervisord]
nodaemon=true

[program:api]
command=python app.py
autostart=true
autorestart=true

[program:worker]
command=python worker.py
autostart=true
autorestart=true
