# supervisord.conf
[supervisord]
nodaemon=true

[program:api]
command=python app.py
autorestart=true

[program:worker]
command=python worker.py
autorestart=true
